# debian/sonic-platform-marvell-dbmvtx9180.postinst
#!/bin/sh
# postinst script for sonic-platform-dbmvtx9180
#
# see: dh_installdeb(1)

depmod -a
systemctl restart kmod
systemctl enable pddf-platform-init.service
systemctl start pddf-platform-init.service
exit 0
